# etc/local.d/20-dmesg-norestrictions.start
#!/bin/sh
if [ -r /proc/sys/kernel/dmesg_restrict ]; then
	if [ $(cat /proc/sys/kernel/dmesg_restrict) = 1 ]; then
		echo 0 > /proc/sys/kernel/dmesg_restrict
	fi
fi
